# ansible/roles/backups/templates/backup.timer.j2
[Unit]
Description=Daily BlindStrader Backup Timer ({{ environment }})

[Timer]
OnCalendar=daily
OnCalendar={{ backup_time }}
Persistent=true

[Install]
WantedBy=timers.target
